# debian/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Reload systemd to register the service (but don't enable it)
        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload || true
            echo "TC358743 setup service installed (disabled by default)"
            echo "To enable: systemctl enable tc358743-setup.service"
            echo "The service should only be enabled after enabling the TC358743 overlay"
        fi
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    *)
        echo "postinst called with unknown argument '$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
